# .config/autostart/easyeffects-service.desktop
[Desktop Entry]
Name=Easy Effects
Exec=easyeffects --gapplication-service
Terminal=false
Type=Application
