# .config/autostart/inhibiter.desktop
[Desktop Entry]
Name=inhibiter
Exec=sh -c "exec ~/.local/lib/inhibiter/inhibiter.py"
Terminal=false
Type=Application
